# sysbackup@.service
[Unit]
Description=Backup the system

[Service]
Type=oneshot
ExecStart=/usr/local/libexec/sysbackup %i
